# templates/upstart.conf.j2
# This file was generated by Ansible for {{ ansible_fqdn }}
# Do NOT modify this file by hand!

{% set srv = item.value %}

description "{{srv.description|default('The service')}}"

start on runlevel [2345]
stop on runlevel [016]

{% for key in srv.env|default({}) %}
env {{key}}={{srv.env[key]}}
{% endfor %}

chdir {{srv.chdir|default(services_chdir)}}
setuid {{srv.user|default(services_user)}}
exec {{srv.exec}}
{{ srv.opts|default('') }}
respawn
console log
